# root/etc/s6-overlay/s6-rc.d/dxhx-svc/run
#!/command/execlineb -P
s6-setuidgid wine
"/usr/bin/dxhx_startup.sh"
